# recipes-graphics/wayland/weston/weston.service.d/wait-drm.conf
[Unit]
# ECU-150v2: wait for the imx-drm display device (card1) before
# starting Weston. On i.MX8MP the DRM minors are:
#   card0 = vivante render node (GPU/ML accelerator, no KMS)
#   card1 = imx-drm display-subsystem (LCDIF3 + HDMI, does modesetting)
# The drm backend in Weston selects card1; we make systemd wait for it.
# card0 is also kept as a Wants= so the render node is up for EGL.
After=dev-dri-card1.device
Wants=dev-dri-card1.device
Wants=dev-dri-card0.device
